# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for stagingserver of SITENAME"

start on net-device-up # 인터넷이 실행될 때 서버를 켠다는 의미
stop on shutdown

respawn # 프로세스가 다운되면 자동으로 재시작한다는 의미

setuid amazingguni # 프로세스 사용자를 amazingguni로 설정
chdir /home/amazingguni/sites/SITENAME/source # 작업 디렉토리 변경

exec ../virtualenv/bin/gunicorn --bind unix:/tmp/SITENAME.socket pomodoro_web.wsgi:application # gunicorn 실행 명령어

